# deploy/supervisor.conf
[program:naruay-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php83 /home/DEPLOY_USER/domains/example.com/laravel/current/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=DEPLOY_USER
numprocs=2
redirect_stderr=true
stdout_logfile=/home/DEPLOY_USER/domains/example.com/laravel/shared/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

[program:naruay-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do /usr/local/bin/php83 /home/DEPLOY_USER/domains/example.com/laravel/current/artisan schedule:run --verbose --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=DEPLOY_USER
redirect_stderr=true
stdout_logfile=/home/DEPLOY_USER/domains/example.com/laravel/shared/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
